# webapp/setup/gunicorn.conf
[program:gunicorn]
command=/usr/bin/gunicorn app:app -b 0.0.0.0:5000
directory=/app