# webapp/setup/supervisor.conf
[program:acworkflow]
command = /home/samk/acworkflow/setup/gunicorn_start.sh
autostart = true
autorestart = true